# docker/resources/etc/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g 'pid /tmp/nginx.pid; daemon off;'